# bin/zenoss
#!/usr/bin/env bash
#
# zenmon        This shell script takes care of starting and stopping zenmon.
#               
# chkconfig: 2345 80 30
# description: Zenoss is a monitoring program written by Zenoss, Inc.

if [ -z "$ZENHOME" ] ; then
  for d in $HOME /home/zenoss ; do 
    for f in .bash_login .bash_profile .profile .bashrc ; do
       if [ -z "$ZENHOME" ] ; then
          [ -f $d/$f ] && . $d/$f
       fi
    done
  done
fi
if [ -z "$ZENHOME" ] ; then
  ZENHOME=/usr/local/zenoss; export ZENHOME
fi

if [ `id -u` -eq 0 ] ; then
   USERNAME=`ls -l $ZENHOME/var/Data.fs | awk '{ print $3 }'`
   if [ "$USERNAME" = "root" -o -z "$USERNAME" ]  ; then 
       echo Unable to determine the right user to run zenoss
       exit 1
   fi
   chown -R $USERNAME $ZENHOME/var
   exec sudo -u $USERNAME $0 $@
fi

. $ZENHOME/bin/zenfunctions
C=
C="$C zeoctl"
C="$C zopectl"
C="$C zenhub"
C="$C zenmodeler"
C="$C zenperfsnmp"
C="$C zenperfxmlrpc"
C="$C zenprocess"
C="$C zenping"
C="$C zensyslog"
C="$C zenstatus"
C="$C zenactions"
C="$C zentrap"
C="$C zencommand"
C="$C zenwin"
C="$C zeneventlog"
CMDS=$C

reverse() {
    args=
    for i
    do
       args="$i $args"
       shift
    done
    echo $args
}

FUNC=$CMD

issue() {
    for cmd in $CMDS
    do
        echo -n "Daemon: $cmd "
        $ZENHOME/bin/$cmd "$FUNC"
    done
}

case "$FUNC" in
  start)
        issue
        ;;
  stop)
        CMDS=`reverse $CMDS`
        issue
        ;;
  restart)
        CMDS=`reverse $CMDS`
	FUNC=stop
	issue
        CMDS=`reverse $CMDS`
	FUNC=start
        issue
        ;;
  status)
        issue
        ;;
  list)
	echo $CMDS
	;;
  *)
        echo $"Usage: $0 {start|stop|restart|status|list}"
        exit 1
esac
exit $?
