# bin/zenoss
#!/usr/bin/env bash
#
# zenmon        This shell script takes care of starting and stopping zenmon.
#               
# chkconfig: 2345 80 30
# description: Zenoss is a monitoring program written by Zenoss, Inc.


. $ZENHOME/bin/zenfunctions
ZOPECMDS="zeoctl zopectl"
ZENCMDS="zenmodeler zenperfsnmp zenprocess zenping zensyslog zenstatus zenactions"

FUNC=$1

issue() {
    for cmd in $CMDS
    do
        echo -n "Daemon: $cmd "
        $ZENHOME/bin/$cmd "$FUNC"
    done
}

case "$1" in
  start)
        CMDS="$ZOPECMDS $ZENCMDS"
        issue
        ;;
  stop)
        CMDS="$ZENCMDS $ZOPECMDS"
        issue
        ;;
  restart)
        CMDS="$ZOPECMDS $ZENCMDS"
        issue
        ;;
  status)
        CMDS="$ZOPECMDS $ZENCMDS"
        issue
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 1
esac
exit $?
